# strings/cfengine/config/linux/etc/cron.weekly/apt-clean
#!/bin/sh
# Auto Clean APT Archives

/usr/bin/apt-get autoclean
